# agent/helpdesk-agent.service
# Usługa systemd dla agenta monitoringu.
# Instalacja: /etc/systemd/system/helpdesk-agent.service
# (prawdziwy klucz API wstawiany przy instalacji — tu placeholder)

[Unit]
Description=Helpdesk monitoring agent
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/helpdesk-agent
Environment=AGENT_API_KEY=CHANGE_ME
ExecStart=/opt/helpdesk-agent/.venv/bin/uvicorn agent:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
